# supervisor.conf
[supervisord]
nodaemon=true

[program:webUI]
command=python /usr/src/app/webUI.py
autostart=true
autorestart=true
stderr_logfile=/var/log/webUI.err.log
stdout_logfile=/var/log/webUI.out.log
